# recipes-jgomez/ejemplo_autostart/files/ejemplo-autostart.service
[Unit]
Description=Hello world auto-startup script

[Service]
ExecStart=/usr/bin/hello-world

[Install]
WantedBy=multi-user.target
